# debian/postinst
#!/bin/sh

update-rc.d callerx defaults s99 k10
/etc/init.d/callerx start